# backups/planck.list
/etc/telegraf/telegraf.conf
/etc/telegraf/telegraf.d/
/etc/systemd/system/telegraf.service.d/
/etc/grafana/grafana.ini
/var/lib/grafana/grafana.db
/etc/prometheus/prometheus.yml
/etc/prometheus/snmp.yml
/etc/prometheus/blackbox.yml
/etc/prometheus/auth.yml
/etc/systemd/system/iperf3.service.d/
