# backups/etc/systemd/system/iperf3.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/iperf3 -s
