# backups/hosts/hadron/configs/latest/etc/systemd/system/unbound.service.d/override.conf
[Unit]
Wants=network-online.target
After=network-online.target systemd-networkd-wait-online.service
